# bmon.socket
[Unit]
Description=BMON gunicorn socket

[Socket]
ListenStream=/run/bmon.sock

[Install]
WantedBy=sockets.target
